# systemd/jarvis-orchestrator.service.d/tls.conf
# Drop-in that switches the orchestrator to HTTPS. Install to
#   /etc/systemd/system/jarvis-orchestrator.service.d/tls.conf  →  systemctl daemon-reload && restart
# Generate the cert/key first with:  bash src/scripts/setup_tls.sh
# Reversible: remove this file + daemon-reload + restart to go back to HTTP.
[Service]
ExecStart=
ExecStart=/usr/local/bin/uv run --no-sync uvicorn main:app --host 0.0.0.0 --port 5000 --workers 1 --ssl-keyfile /srv/jarvis/tls/server.key --ssl-certfile /srv/jarvis/tls/server.crt
